# etc/cron.monthly/scrollkeeper
#!/bin/sh

set -e

[ -x /usr/bin/scrollkeeper-rebuilddb ] || exit 0

umask 022
scrollkeeper-rebuilddb -q

