# etc/cron.daily/vivipos
#!/bin/sh

# check daily backup 
bak=/data/backups
bak_scripts=/data/scripts/backup.sh

if ([ -d $bak ] && [ -f $bak_scripts ]); then
  exec $bak_scripts;
fi
